# provision/playbooks/templates/gunicorn/gunicorn.socket.j2
[Unit]
Description=gunicorn socket

[Socket]
ListenStream="{{ aws.instance.gunicorn.socket.file.path }}"

SocketUser=nginx


[Install]
WantedBy=sockets.target